# package/files/ubuntu/lib/systemd/system/ribd.service
[Unit]
Description=ribd
After=network.target

[Service]
Type=simple
User=root
Group=root
StandardOutput=syslog
StandardError=syslog
ExecStart=/usr/bin/ribd
ExecStop=/bin/kill -INT ${MAINPID}
Restart=always

[Install]
WantedBy=multi-user.target
